# noip/noip.timer
[Unit]
Description=No-IP Dynamic DNS Update Client

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
